# update-motd.d/00-header
#!/bin/sh
[ -r /etc/os-release ] && . /etc/os-release
printf "\033[1m-> Welcome to %s\033[m\n" "$PRETTY_NAME"
printf 'Message last updated at %s\n' "$(date)"
