# update-motd.d/10-help-text
#!/bin/sh
printf "\n\033[1m-> Help\033[m\n"
printf " - Documentation: https://docs.voidlinux.org/\n"
